# tms_rc/tms_rc_smartpal/build-smartpal_control-Desktop-Default/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/rt/catkin_ws/src/ros_tms/tms_rc/tms_rc_smartpal/build-smartpal_control-Desktop-Default/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
